# config/files/usr/etc/profile.d/settheming.sh
gsettings set org.gnome.desktop.interface gtk-theme 'Catppuccin-Mocha-Standard-Green-Dark'
gsettings set org.gnome.desktop.interface color-scheme 'prefer-dark'
gsettings set org.gnome.desktop.interface icon-theme 'ePapirus'
